# systemd/rss-to-kobo@.service
[Unit]
Description=RSS to Kobo for profile %i
After=network.target

[Service]
Type=oneshot
User=YOUR_USERNAME  # Change this to the user that will run the service
WorkingDirectory=/path/to/rss-to-kobo
Environment=PYTHONUNBUFFERED=1

# Use the Python from your virtual environment
ExecStart=/path/to/venv/bin/python rss_to_kobo.py --user %i

# Logging configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rss-to-kobo-%i

# Security options
ProtectSystem=full
PrivateTmp=yes
NoNewPrivileges=yes

# Restart on failure (but not for token expiration)
Restart=on-failure
RestartSec=60s

[Install]
WantedBy=multi-user.target
